# modeco/configs/ingest_journalctl/journal.output-to-log.service
[Unit]
After=

[Service]
ExecStart=/infometis/journalctl/journalctl-output-to-log

[Install]
WantedBy=default.target
